# rag/launchd/com.dotfiles.rag-indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles.rag-indexer</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.rag/venv/bin/python</string>
        <string>-m</string>
        <string>src.indexer</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__HOME__/Developer/dotfiles/rag</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.rag/logs/indexer.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.rag/logs/indexer.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
