# scripts/com.dotfiles.weekly-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles.weekly-cleanup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/Developer/dotfiles/scripts/weekly-cleanup.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>604800</integer>
    <key>StandardOutPath</key>
    <string>/tmp/weekly-cleanup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/weekly-cleanup.log</string>
</dict>
</plist>
